# mitre-attack/mappings/rule_technique_mapping.json
{
  "version": "14.0",
  "lab": "Enterprise Detection Engineering Lab",
  "author": "Chandra Sekhar Chakraborty",
  "last_updated": "2026-01-15",
  "mappings": [
    {
      "rule_id": "a3c2f1e8-9b4d-4c7a-8e2f-1d6b3a9c5f0e",
      "rule_name": "Mimikatz LSASS Memory Access",
      "tactic": "Credential Access",
      "tactic_id": "TA0006",
      "technique": "OS Credential Dumping: LSASS Memory",
      "technique_id": "T1003.001",
      "platform": ["Windows"],
      "data_source": "Process Access",
      "log_source": "Sysmon Event ID 10",
      "severity": "critical",
      "validation_status": "validated",
      "atomic_test_id": "T1003.001-1",
      "sigma_file": "detection-rules/sigma/credential_access_mimikatz_lsass.yml"
    },
    {
      "rule_id": "b7d3e2a1-4f8c-4e9b-9c1d-2a7f5b3d8e4c",
      "rule_name": "PowerShell Obfuscated Script Execution",
      "tactic": "Defense Evasion",
      "tactic_id": "TA0005",
      "technique": "Obfuscated Files or Information: Command Obfuscation",
      "technique_id": "T1027.010",
      "platform": ["Windows"],
      "data_source": "Process Creation",
      "log_source": "Sysmon Event ID 1",
      "severity": "high",
      "validation_status": "validated",
      "atomic_test_id": "T1027.010-1",
      "sigma_file": "detection-rules/sigma/defense_evasion_powershell_obfuscation.yml"
    },
    {
      "rule_id": "c9e4f3b2-5a1d-4f8e-8b2c-3d9a7c4f1e6b",
      "rule_name": "PsExec Lateral Movement via Network Share",
      "tactic": "Lateral Movement",
      "tactic_id": "TA0008",
      "technique": "Remote Services: SMB/Windows Admin Shares",
      "technique_id": "T1021.002",
      "platform": ["Windows"],
      "data_source": "Process Creation, Windows Service",
      "log_source": "Sysmon Event ID 1 / System Event 7045",
      "severity": "high",
      "validation_status": "validated",
      "atomic_test_id": "T1021.002-1",
      "sigma_file": "detection-rules/sigma/lateral_movement_psexec.yml"
    },
    {
      "rule_id": "d1f5a4c3-6b2e-4a9f-7c3d-4e8b1a5f2d9c",
      "rule_name": "WMI Spawning Suspicious Child Process",
      "tactic": "Execution",
      "tactic_id": "TA0002",
      "technique": "Windows Management Instrumentation",
      "technique_id": "T1047",
      "platform": ["Windows"],
      "data_source": "Process Creation",
      "log_source": "Sysmon Event ID 1",
      "severity": "high",
      "validation_status": "validated",
      "atomic_test_id": "T1047-1",
      "sigma_file": "detection-rules/sigma/execution_wmi_process_spawn.yml"
    },
    {
      "rule_id": "e2a6b5d4-7c3f-4b0a-8d4e-5f9c2b6a3e0d",
      "rule_name": "Suspicious Scheduled Task Creation via CLI",
      "tactic": "Persistence",
      "tactic_id": "TA0003",
      "technique": "Scheduled Task/Job: Scheduled Task",
      "technique_id": "T1053.005",
      "platform": ["Windows"],
      "data_source": "Process Creation, Scheduled Job",
      "log_source": "Sysmon Event ID 1",
      "severity": "medium",
      "validation_status": "validated",
      "atomic_test_id": "T1053.005-1",
      "sigma_file": "detection-rules/sigma/persistence_scheduled_task.yml"
    },
    {
      "rule_id": "f3b7c6e5-8d4a-4c1b-9e5f-6a0d3c7b4f1e",
      "rule_name": "Registry Run Key Persistence Modification",
      "tactic": "Persistence",
      "tactic_id": "TA0003",
      "technique": "Boot or Logon Autostart Execution: Registry Run Keys",
      "technique_id": "T1547.001",
      "platform": ["Windows"],
      "data_source": "Windows Registry",
      "log_source": "Sysmon Event ID 13",
      "severity": "medium",
      "validation_status": "validated",
      "atomic_test_id": "T1547.001-1",
      "sigma_file": "detection-rules/sigma/persistence_registry_run_key.yml"
    },
    {
      "rule_id": "a4c8d7f6-9e5b-4d2c-0f6a-7b1e4d8c5a2f",
      "rule_name": "Account Discovery via Net User Commands",
      "tactic": "Discovery",
      "tactic_id": "TA0007",
      "technique": "Account Discovery: Local Account",
      "technique_id": "T1087.001",
      "platform": ["Windows"],
      "data_source": "Process Creation",
      "log_source": "Sysmon Event ID 1",
      "severity": "low",
      "validation_status": "validated",
      "atomic_test_id": "T1087.001-1",
      "sigma_file": "detection-rules/sigma/discovery_net_user_enum.yml"
    },
    {
      "rule_id": "b5d9e8a7-0f6c-4e3d-1a7b-8c2f5e9d6b3a",
      "rule_name": "Pass-the-Hash Attack via NTLM Authentication",
      "tactic": "Lateral Movement",
      "tactic_id": "TA0008",
      "technique": "Use Alternate Authentication Material: Pass the Hash",
      "technique_id": "T1550.002",
      "platform": ["Windows"],
      "data_source": "Logon Session, User Account Authentication",
      "log_source": "Windows Security Event ID 4624",
      "severity": "high",
      "validation_status": "validated",
      "atomic_test_id": "T1550.002-1",
      "sigma_file": "detection-rules/sigma/lateral_movement_pass_the_hash.yml"
    },
    {
      "rule_id": "c6e0f9b8-1a7d-4f4e-2b8c-9d3a6f0e7c4b",
      "rule_name": "Suspicious Encoded Command Execution",
      "tactic": "Defense Evasion",
      "tactic_id": "TA0005",
      "technique": "Deobfuscate/Decode Files or Information",
      "technique_id": "T1140",
      "platform": ["Windows"],
      "data_source": "Process Creation",
      "log_source": "Sysmon Event ID 1",
      "severity": "high",
      "validation_status": "validated",
      "atomic_test_id": "T1140-1",
      "sigma_file": "detection-rules/sigma/defense_evasion_encoded_command.yml"
    },
    {
      "rule_id": "d7f1a0c9-2b8e-4a5f-3c9d-0e4b7a1f8d5c",
      "rule_name": "LOLBAS Abuse via Certutil Download or Decode",
      "tactic": "Defense Evasion",
      "tactic_id": "TA0005",
      "technique": "System Binary Proxy Execution",
      "technique_id": "T1218",
      "platform": ["Windows"],
      "data_source": "Process Creation, Network Traffic",
      "log_source": "Sysmon Event ID 1, 3",
      "severity": "high",
      "validation_status": "validated",
      "atomic_test_id": "T1218-14",
      "sigma_file": "detection-rules/sigma/defense_evasion_certutil_lolbas.yml"
    }
  ]
}
